# supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=python manage.py migrate --noinput
command=python manage.py collectstatic --noinput
command=python manage.py runserver 0.0.0.0:8000
directory=/code
autostart=true
autorestart=true
stderr_logfile=/var/log/django.err.log
stdout_logfile=/var/log/django.out.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
